# nginx/etc/periodic/daily/nginx
#!/bin/sh
echo "reloading nginx" >> /var/log/cron.log 2>&1
nginx reload >> /var/log/cron.log 2>&1
