# certbot/etc/cron.d/certbot
23 3 * * * root /opt/certbot/venv/bin/certbot renew >> /var/log/cron.log 2>&1
